# docker/php/crontab
*/5 * * * * php /path/to/yii session/update-completed >> /var/log/cron.log 2>&1
